# woabot.service
[Unit]
Description=Woabot
After=multi-user.target

[Service]
ExecStart=/home/daniel/woabot/woabot

[Install]
WantedBy=multi-user.target